# shadowsocks.conf
[program:shadowsocks]
command=ssserver -c /etc/shadowsocks.json -d
autorestart=true  
user=nobody 
